# roles/nfs_client/templates/nfs.automount.j2
[Unit]
Description=Automount NFS {{ item.path }}
ConditionPathExists={{ item.path }}

[Automount]
Where={{ item.path }}

[Install]
WantedBy=multi-user.target
